# config.example/supervisord.example.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/bin/bash apache2-foreground
